# images/promhealth/s6/services.d/exporter/run
#!/usr/bin/with-contenv sh

exec /usr/bin/node /usr/local/src/exporter/exporter.js
