# prometheus/supervisord.conf
[supervisord]
nodaemon = true
loglevel = info
logfile=/dev/null
logfile_maxbytes=0
pidfile  = /run/supervisord.pid
logfile_maxbytes = 20MB
logfile_backups = 2

[unix_http_server]
file = /run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///run/supervisor.sock

[program:dockergen]
command = docker-gen -watch -notify "/usr/local/bin/reload-prom.sh" /etc/prometheus/prometheus.yaml.tmpl /etc/prometheus/generated.yaml
numprocs = 1
process_name = %(program_name)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart = true
autostart = true

[program:prometheus]
command = /usr/local/bin/prometheus-wrapper.sh
numprocs = 1
process_name = %(program_name)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart = true
autostart = true
