# nextcloud-nginx/s6/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

exec /entrypoint.sh php-fpm