# images/dnsmasq/s6/services.d/dnsmasq/run
#!/command/execlineb -P

dnsmasq --config-file /etc/dnsmaq.conf
